# qetesh/qetesh.service.in
[Unit]
Description=Qetesh web application server

[Service]
ExecStart=@prefix@/bin/qetesh
Restart=always
IgnoreSIGPIPE=no
StandardOutput=syslog
StandardError=inherit

[Install]
Alias=qetesh.service
